# deploy/cariar-bot.service
[Unit]
Description=CarIA — WhatsApp Bot (Garcia Multimarcas)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cariar
EnvironmentFile=/opt/cariar/.env
ExecStart=/opt/cariar/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
